# sigint-ops.service
[Unit]
Description=SIGINT Ops Backend
After=network.target

[Service]
User=ubuntu
WorkingDirectory=/home/ubuntu/sigint-ops/backend
ExecStart=/home/ubuntu/sigint-ops/backend/venv/bin/python app.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
